# package/base-files/files/etc/uci-defaults/99-ser2net
uci -q batch << EOI
set ser2net.@proxy[1].enabled='1'
set ser2net.@proxy[1].protocol='raw'
commit ser2net
EOI
